# launchd/com.orchestrate.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<!--
  task-orchestrate monitor web UI — launchd-managed, auto-restart on crash/code change.

  BEFORE INSTALLING — set PROJECT_DIR and run install-launchd.sh:

    PROJECT_DIR="$HOME/apps/ai-console" bash skills/task-orchestrate/bin/install-launchd.sh

  Logs (relative to WorkingDirectory):
    stdout -> .orchestrate/logs/monitor.log
    stderr -> .orchestrate/logs/monitor-err.log
-->

<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.orchestrate.monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BIN__</string>
    <string>__MONITOR_SERVER__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>~/apps/my-project</string>

  <key>StandardOutPath</key>
  <string>~/apps/my-project/.orchestrate/logs/monitor.log</string>

  <key>StandardErrorPath</key>
  <string>~/apps/my-project/.orchestrate/logs/monitor-err.log</string>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

</dict>
</plist>
